# Lab4/src/task/part7/mymsg.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/mymsg.sh
RemainAfterExit=yes