# Lab4/src/task/part5/task2.mount
[Unit]
Description=Mounting a partition /mnt/mydata

[Mount]
What=/dev/sdb1
Where=/mnt/mydata
Type=ext4
Options=defaults

[Install]
WantedBy=multi-user.target

