# Lab4/src/task/part6/task1.automount
[Unit]
Description=Automount for /mnt/mydata

[Automount]
Where=/mnt/mydata
TimeoutIdleSec=30

[Install]
WantedBy=multi-user.target
